# backup/profile.d/00-config.sh
# shellcheck shell=sh

export DIRENV_CONFIG="${RC_CONFIG}/direnv/config.toml"
export INPUTRC="${RC_CONFIG}/readline/inputrc"
export STARSHIP_CONFIG="${RC_CONFIG}/config.toml"
