# rc/profile.d/00-constants.sh
# shellcheck shell=sh

export HISTAPPEND="history -a; history -c; history -r"
export TERMINAL_TITLE="\[\e]0;\h@\u: \w\a\]"
