# rc/profile.d/00-config.sh
# shellcheck shell=sh

export GIT_CONFIG_SYSTEM="/etc/gitconfig"
export GIT_TEMPLATE_DIR="${RC_CONFIG}/git/template"
export INPUTRC="${RC_CONFIG}/readline/inputrc"
export STARSHIP_CONFIG="${RC_CONFIG}/config.toml"
