# rc/profile.d/00-constants.sh
# shellcheck shell=sh

export PROMPT_COMMAND="history -a; history -c; history -r"

#export CHARSET=UTF-8
#export LANG=C.UTF-8
#export LC_COLLATE=C
