# rc/profile.d/symbol.sh
# shellcheck shell=sh

#
# Color Symbols library (generated by: color)

# Green "✔"
#
export GreenOk="\e[1;32m✔\e[0m"

# Green "✔" Escaped
#
export GreenOkEsc="\[\e[1;32m\]✔\e[0m"

# Red "✘"
#
export RedX="\e[1;32m✘\e[0m"

# Red "✘" Escaped
#
export RedXEsc="\[\e[1;32m\]✘\e[0m"


# YellowMark "!"
#
export YellowMark="\e[1;33m!\e[0m"

# YellowMark "✘" Escaped
#
export YellowMarkEsc="\[\e[1;33m\]!\e[0m"
