# rc/profile.d/symbol.sh
# shellcheck shell=sh

#
# Symbol Colors library (generated by: color)

# <html><h2>Critical Icon</h2>
# <p><strong><code>$Critical</code></strong> Critical Icon.</p>
# </html>
export CriticalIcon="✘"

# <html><h2>Critical </h2>
# <p><strong><code>$Critical</code></strong> Critical  RedInvert.</p>
# </html>
export Critical="\033[7;31m✘\033[0m"

# <html><h2>Critical Esc</h2>
# <p><strong><code>$CriticalEsc</code></strong> Critical Esc RedInvert.</p>
# </html>
export CriticalEsc="\[\033[7;31m✘\033[0m\]"

# <html><h2>Critical Zsh</h2>
# <p><strong><code>$CriticalZsh</code></strong> Critical Zsh RedInvert.</p>
# </html>
export CriticalZsh="%{\033[7;31m✘\033[0m%}"

# <html><h2>Cube Icon</h2>
# <p><strong><code>$Cube</code></strong> Cube Icon.</p>
# </html>
export CubeIcon="❒"

# <html><h2>Cube </h2>
# <p><strong><code>$Cube</code></strong> Cube  Cyan.</p>
# </html>
export Cube="\033[1;36m❒\033[0m"

# <html><h2>Cube Esc</h2>
# <p><strong><code>$CubeEsc</code></strong> Cube Esc Cyan.</p>
# </html>
export CubeEsc="\[\033[1;36m❒\033[0m\]"

# <html><h2>Cube Zsh</h2>
# <p><strong><code>$CubeZsh</code></strong> Cube Zsh Cyan.</p>
# </html>
export CubeZsh="%{\033[1;36m❒\033[0m%}"

# <html><h2>Debug Icon</h2>
# <p><strong><code>$Debug</code></strong> Debug Icon.</p>
# </html>
export DebugIcon="＋"

# <html><h2>Debug </h2>
# <p><strong><code>$Debug</code></strong> Debug  Magenta.</p>
# </html>
export Debug="\033[1;35m＋\033[0m"

# <html><h2>Debug Esc</h2>
# <p><strong><code>$DebugEsc</code></strong> Debug Esc Magenta.</p>
# </html>
export DebugEsc="\[\033[1;35m＋\033[0m\]"

# <html><h2>Debug Zsh</h2>
# <p><strong><code>$DebugZsh</code></strong> Debug Zsh Magenta.</p>
# </html>
export DebugZsh="%{\033[1;35m＋\033[0m%}"

# <html><h2>Division Icon</h2>
# <p><strong><code>$Division</code></strong> Division Icon.</p>
# </html>
export DivisionIcon="∕"

# <html><h2>Division </h2>
# <p><strong><code>$Division</code></strong> Division  Yellow.</p>
# </html>
export Division="\033[1;33m∕\033[0m"

# <html><h2>Division Esc</h2>
# <p><strong><code>$DivisionEsc</code></strong> Division Esc Yellow.</p>
# </html>
export DivisionEsc="\[\033[1;33m∕\033[0m\]"

# <html><h2>Division Zsh</h2>
# <p><strong><code>$DivisionZsh</code></strong> Division Zsh Yellow.</p>
# </html>
export DivisionZsh="%{\033[1;33m∕\033[0m%}"

# <html><h2>Ellipsis Icon</h2>
# <p><strong><code>$Ellipsis</code></strong> Ellipsis Icon.</p>
# </html>
export EllipsisIcon="…"

# <html><h2>Ellipsis </h2>
# <p><strong><code>$Ellipsis</code></strong> Ellipsis  Cyan.</p>
# </html>
export Ellipsis="\033[1;36m…\033[0m"

# <html><h2>Ellipsis Esc</h2>
# <p><strong><code>$EllipsisEsc</code></strong> Ellipsis Esc Cyan.</p>
# </html>
export EllipsisEsc="\[\033[1;36m…\033[0m\]"

# <html><h2>Ellipsis Zsh</h2>
# <p><strong><code>$EllipsisZsh</code></strong> Ellipsis Zsh Cyan.</p>
# </html>
export EllipsisZsh="%{\033[1;36m…\033[0m%}"

# <html><h2>Error Icon</h2>
# <p><strong><code>$Error</code></strong> Error Icon.</p>
# </html>
export ErrorIcon="✘"

# <html><h2>Error </h2>
# <p><strong><code>$Error</code></strong> Error  Red.</p>
# </html>
export Error="\033[1;31m✘\033[0m"

# <html><h2>Error Esc</h2>
# <p><strong><code>$ErrorEsc</code></strong> Error Esc Red.</p>
# </html>
export ErrorEsc="\[\033[1;31m✘\033[0m\]"

# <html><h2>Error Zsh</h2>
# <p><strong><code>$ErrorZsh</code></strong> Error Zsh Red.</p>
# </html>
export ErrorZsh="%{\033[1;31m✘\033[0m%}"

# <html><h2>Finish Icon</h2>
# <p><strong><code>$Finish</code></strong> Finish Icon.</p>
# </html>
export FinishIcon="＞"

# <html><h2>Finish </h2>
# <p><strong><code>$Finish</code></strong> Finish  Green.</p>
# </html>
export Finish="=\033[1;32m＞\033[0m"

# <html><h2>Finish Esc</h2>
# <p><strong><code>$FinishEsc</code></strong> Finish Esc Green.</p>
# </html>
export FinishEsc="\[=\033[1;32m＞\033[0m\]"

# <html><h2>Finish Zsh</h2>
# <p><strong><code>$FinishZsh</code></strong> Finish Zsh Green.</p>
# </html>
export FinishZsh="%{=\033[1;32m＞\033[0m%}"

# <html><h2>Harpoon Icon</h2>
# <p><strong><code>$Harpoon</code></strong> Harpoon Icon.</p>
# </html>
export HarpoonIcon="⇌"

# <html><h2>Harpoon </h2>
# <p><strong><code>$Harpoon</code></strong> Harpoon  Cyan.</p>
# </html>
export Harpoon="\033[1;36m⇌\033[0m"

# <html><h2>Harpoon Esc</h2>
# <p><strong><code>$HarpoonEsc</code></strong> Harpoon Esc Cyan.</p>
# </html>
export HarpoonEsc="\[\033[1;36m⇌\033[0m\]"

# <html><h2>Harpoon Zsh</h2>
# <p><strong><code>$HarpoonZsh</code></strong> Harpoon Zsh Cyan.</p>
# </html>
export HarpoonZsh="%{\033[1;36m⇌\033[0m%}"

# <html><h2>Left Icon</h2>
# <p><strong><code>$Left</code></strong> Left Icon.</p>
# </html>
export LeftIcon="〈"

# <html><h2>Left </h2>
# <p><strong><code>$Left</code></strong> Left  Green.</p>
# </html>
export Left="\033[1;32m〈\033[0m"

# <html><h2>Left Esc</h2>
# <p><strong><code>$LeftEsc</code></strong> Left Esc Green.</p>
# </html>
export LeftEsc="\[\033[1;32m〈\033[0m\]"

# <html><h2>Left Zsh</h2>
# <p><strong><code>$LeftZsh</code></strong> Left Zsh Green.</p>
# </html>
export LeftZsh="%{\033[1;32m〈\033[0m%}"

# <html><h2>Minus Icon</h2>
# <p><strong><code>$Minus</code></strong> Minus Icon.</p>
# </html>
export MinusIcon="-"

# <html><h2>Minus </h2>
# <p><strong><code>$Minus</code></strong> Minus  Red.</p>
# </html>
export Minus="\033[1;31m-\033[0m"

# <html><h2>Minus Esc</h2>
# <p><strong><code>$MinusEsc</code></strong> Minus Esc Red.</p>
# </html>
export MinusEsc="\[\033[1;31m-\033[0m\]"

# <html><h2>Minus Zsh</h2>
# <p><strong><code>$MinusZsh</code></strong> Minus Zsh Red.</p>
# </html>
export MinusZsh="%{\033[1;31m-\033[0m%}"

# <html><h2>More Icon</h2>
# <p><strong><code>$More</code></strong> More Icon.</p>
# </html>
export MoreIcon="＞"

# <html><h2>More </h2>
# <p><strong><code>$More</code></strong> More  Magenta.</p>
# </html>
export More="\033[1;35m＞\033[0m"

# <html><h2>More Esc</h2>
# <p><strong><code>$MoreEsc</code></strong> More Esc Magenta.</p>
# </html>
export MoreEsc="\[\033[1;35m＞\033[0m\]"

# <html><h2>More Zsh</h2>
# <p><strong><code>$MoreZsh</code></strong> More Zsh Magenta.</p>
# </html>
export MoreZsh="%{\033[1;35m＞\033[0m%}"

# <html><h2>Multiply Icon</h2>
# <p><strong><code>$Multiply</code></strong> Multiply Icon.</p>
# </html>
export MultiplyIcon="×"

# <html><h2>Multiply </h2>
# <p><strong><code>$Multiply</code></strong> Multiply  Blue.</p>
# </html>
export Multiply="\033[1;34m×\033[0m"

# <html><h2>Multiply Esc</h2>
# <p><strong><code>$MultiplyEsc</code></strong> Multiply Esc Blue.</p>
# </html>
export MultiplyEsc="\[\033[1;34m×\033[0m\]"

# <html><h2>Multiply Zsh</h2>
# <p><strong><code>$MultiplyZsh</code></strong> Multiply Zsh Blue.</p>
# </html>
export MultiplyZsh="%{\033[1;34m×\033[0m%}"

# <html><h2>Notice Icon</h2>
# <p><strong><code>$Notice</code></strong> Notice Icon.</p>
# </html>
export NoticeIcon="！‼"

# <html><h2>Notice </h2>
# <p><strong><code>$Notice</code></strong> Notice  Yellow.</p>
# </html>
export Notice="\033[1;33m！‼\033[0m"

# <html><h2>Notice Esc</h2>
# <p><strong><code>$NoticeEsc</code></strong> Notice Esc Yellow.</p>
# </html>
export NoticeEsc="\[\033[1;33m！‼\033[0m\]"

# <html><h2>Notice Zsh</h2>
# <p><strong><code>$NoticeZsh</code></strong> Notice Zsh Yellow.</p>
# </html>
export NoticeZsh="%{\033[1;33m！‼\033[0m%}"

# <html><h2>Ok Icon</h2>
# <p><strong><code>$Ok</code></strong> Ok Icon.</p>
# </html>
export OkIcon="✔"

# <html><h2>Ok </h2>
# <p><strong><code>$Ok</code></strong> Ok  Green.</p>
# </html>
export Ok="\033[1;32m✔\033[0m"

# <html><h2>Ok Esc</h2>
# <p><strong><code>$OkEsc</code></strong> Ok Esc Green.</p>
# </html>
export OkEsc="\[\033[1;32m✔\033[0m\]"

# <html><h2>Ok Zsh</h2>
# <p><strong><code>$OkZsh</code></strong> Ok Zsh Green.</p>
# </html>
export OkZsh="%{\033[1;32m✔\033[0m%}"

# <html><h2>Plus Icon</h2>
# <p><strong><code>$Plus</code></strong> Plus Icon.</p>
# </html>
export PlusIcon="+"

# <html><h2>Plus </h2>
# <p><strong><code>$Plus</code></strong> Plus  Green.</p>
# </html>
export Plus="\033[1;32m+\033[0m"

# <html><h2>Plus Esc</h2>
# <p><strong><code>$PlusEsc</code></strong> Plus Esc Green.</p>
# </html>
export PlusEsc="\[\033[1;32m+\033[0m\]"

# <html><h2>Plus Zsh</h2>
# <p><strong><code>$PlusZsh</code></strong> Plus Zsh Green.</p>
# </html>
export PlusZsh="%{\033[1;32m+\033[0m%}"

# <html><h2>Right Icon</h2>
# <p><strong><code>$Right</code></strong> Right Icon.</p>
# </html>
export RightIcon="〉"

# <html><h2>Right </h2>
# <p><strong><code>$Right</code></strong> Right  Green.</p>
# </html>
export Right="\033[1;32m〉\033[0m"

# <html><h2>Right Esc</h2>
# <p><strong><code>$RightEsc</code></strong> Right Esc Green.</p>
# </html>
export RightEsc="\[\033[1;32m〉\033[0m\]"

# <html><h2>Right Zsh</h2>
# <p><strong><code>$RightZsh</code></strong> Right Zsh Green.</p>
# </html>
export RightZsh="%{\033[1;32m〉\033[0m%}"

# <html><h2>Start Icon</h2>
# <p><strong><code>$Start</code></strong> Start Icon.</p>
# </html>
export StartIcon="＞"

# <html><h2>Start </h2>
# <p><strong><code>$Start</code></strong> Start  Blue.</p>
# </html>
export Start="=\033[1;34m＞\033[0m"

# <html><h2>Start Esc</h2>
# <p><strong><code>$StartEsc</code></strong> Start Esc Blue.</p>
# </html>
export StartEsc="\[=\033[1;34m＞\033[0m\]"

# <html><h2>Start Zsh</h2>
# <p><strong><code>$StartZsh</code></strong> Start Zsh Blue.</p>
# </html>
export StartZsh="%{=\033[1;34m＞\033[0m%}"

# <html><h2>Success Icon</h2>
# <p><strong><code>$Success</code></strong> Success Icon.</p>
# </html>
export SuccessIcon="◉"

# <html><h2>Success </h2>
# <p><strong><code>$Success</code></strong> Success  Blue.</p>
# </html>
export Success="\033[1;34m◉\033[0m"

# <html><h2>Success Esc</h2>
# <p><strong><code>$SuccessEsc</code></strong> Success Esc Blue.</p>
# </html>
export SuccessEsc="\[\033[1;34m◉\033[0m\]"

# <html><h2>Success Zsh</h2>
# <p><strong><code>$SuccessZsh</code></strong> Success Zsh Blue.</p>
# </html>
export SuccessZsh="%{\033[1;34m◉\033[0m%}"

# <html><h2>Tilde Icon</h2>
# <p><strong><code>$Tilde</code></strong> Tilde Icon.</p>
# </html>
export TildeIcon="～"

# <html><h2>Tilde </h2>
# <p><strong><code>$Tilde</code></strong> Tilde  Blue.</p>
# </html>
export Tilde="\033[1;34m～\033[0m"

# <html><h2>Tilde Esc</h2>
# <p><strong><code>$TildeEsc</code></strong> Tilde Esc Blue.</p>
# </html>
export TildeEsc="\[\033[1;34m～\033[0m\]"

# <html><h2>Tilde Zsh</h2>
# <p><strong><code>$TildeZsh</code></strong> Tilde Zsh Blue.</p>
# </html>
export TildeZsh="%{\033[1;34m～\033[0m%}"

# <html><h2>Verbose Icon</h2>
# <p><strong><code>$Verbose</code></strong> Verbose Icon.</p>
# </html>
export VerboseIcon="＞"

# <html><h2>Verbose </h2>
# <p><strong><code>$Verbose</code></strong> Verbose  Cyan.</p>
# </html>
export Verbose="\033[1;36m＞\033[0m"

# <html><h2>Verbose Esc</h2>
# <p><strong><code>$VerboseEsc</code></strong> Verbose Esc Cyan.</p>
# </html>
export VerboseEsc="\[\033[1;36m＞\033[0m\]"

# <html><h2>Verbose Zsh</h2>
# <p><strong><code>$VerboseZsh</code></strong> Verbose Zsh Cyan.</p>
# </html>
export VerboseZsh="%{\033[1;36m＞\033[0m%}"

# <html><h2>Warning Icon</h2>
# <p><strong><code>$Warning</code></strong> Warning Icon.</p>
# </html>
export WarningIcon="‼"

# <html><h2>Warning </h2>
# <p><strong><code>$Warning</code></strong> Warning  Yellow.</p>
# </html>
export Warning="\033[1;33m‼\033[0m"

# <html><h2>Warning Esc</h2>
# <p><strong><code>$WarningEsc</code></strong> Warning Esc Yellow.</p>
# </html>
export WarningEsc="\[\033[1;33m‼\033[0m\]"

# <html><h2>Warning Zsh</h2>
# <p><strong><code>$WarningZsh</code></strong> Warning Zsh Yellow.</p>
# </html>
export WarningZsh="%{\033[1;33m‼\033[0m%}"

