# etc/profile.d/00-constants.sh
# shellcheck shell=sh

export HISTAPPEND="history -a; history -c; history -r"
export TERMINAL_TITLE="\[\e]0;\h@\u: \w\a\]"

#export CHARSET=UTF-8
#export LANG=C.UTF-8
#export LC_COLLATE=C
